# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:openclaw-gateway]
; 前台运行，supervisord 管理生命周期
command=openclaw gateway run --port 18789 --allow-unconfigured
directory=/root
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
priority=10

[program:clawpanel]
command=node /app/scripts/serve.js
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
; 等 Gateway 启动后再启动面板
startsecs=5
priority=20